# logseq.conf
[program:logseq]
user=logseq
environment=HOME="/home/logseq",USER="logseq"
command=/Logseq-linux-x64/Logseq --no-sandbox
autorestart=true
